# launchd/com.telepi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.telepi</string>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/TelePi</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/node</string>
    <string>/ABSOLUTE/PATH/TO/TelePi/dist/index.js</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/telepi.out.log</string>

  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/telepi.err.log</string>
</dict>
</plist>
